# projects/infrastructure/pi3/etc/cron.d/camera
SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

# camera
*/5 * * * *   root    bash /home/pi/projects/infrastructure-raspberrypi/monitoring/security-cam.sh